# dots/Library/LaunchAgents/com.rizowski.brew-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.rizowski.brew-check</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-lc</string>
    <string>mkdir -p "$HOME/.cache/dot-brew" &amp;&amp; exec "$HOME/bin/brew-check" &gt;&gt;"$HOME/.cache/dot-brew/check.log" 2&gt;&amp;1</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>86400</integer>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
